# nginx-tomcat/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:tomcat]
command=catalina.sh run

